# linux/roles/scu/files/docker-login-refresh.timer
[Unit]
Description=Docker login timer

[Timer]
OnCalendar=00/6:00

[Install]
WantedBy=timers.target
