# linux/roles/scu/files/docker-login-refresh.timer
[Unit]
Description=Docker login timer

[Timer]
OnActiveSec=21600

[Install]
WantedBy=timers.target
