# linux/roles/scu/files/docker-login-refresh.service
[Unit]
Description=Docker login refresh

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart docker-login.service
